# systemd/aura-sensor.service
[Unit]
Description=Aura heartbeat sensor
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/apps/aura
ExecStart=/usr/bin/python3 /opt/apps/aura/sensor/aura_sensor.py run
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
